# servers/intel/intel-daily.timer
[Unit]
Description=Run daily tech intel batch at 07:00

[Timer]
OnCalendar=*-*-* 07:00:00
Persistent=true

[Install]
WantedBy=timers.target
