# servers/intel/intel-daily.timer
[Unit]
Description=Run daily tech intel batch at 04:00 KST

[Timer]
OnCalendar=*-*-* 04:00:00 Asia/Seoul
Persistent=true

[Install]
WantedBy=timers.target
